# installer/systemd/tranzfer-edi.service
[Unit]
Description=TranzFer MFT — edi-converter
After=network-online.target
Wants=network-online.target
Documentation=https://docs.tranzfer.io
# Stateless — no database dependency

[Service]
User=tranzfer
Group=tranzfer
WorkingDirectory=/opt/tranzfer-mft
ExecStart=/usr/bin/java \
    -Xms256m -Xmx512m \
    -XX:+UseZGC -XX:+ZGenerational \
    -XX:+EnableDynamicAgentLoading \
    --add-opens java.base/java.lang=ALL-UNNAMED \
    --add-opens java.base/java.lang.reflect=ALL-UNNAMED \
    --add-opens java.base/java.util=ALL-UNNAMED \
    --add-opens java.base/java.lang.invoke=ALL-UNNAMED \
    --add-opens java.base/java.io=ALL-UNNAMED \
    -Dserver.port=8095 \
    -Dspring.profiles.active=production \
    -Dspring.config.additional-location=/opt/tranzfer-mft/config/ \
    -jar /opt/tranzfer-mft/lib/edi-converter.jar
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tranzfer-edi

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/tranzfer-mft /tmp
PrivateTmp=true

[Install]
WantedBy=multi-user.target
