# installer/systemd/tranzfer-ftp.service
[Unit]
Description=TranzFer MFT — ftp-service
After=network-online.target postgresql.service redis.service rabbitmq-server.service
Wants=network-online.target
Documentation=https://docs.tranzfer.io
# Also listens on FTP port 21 (requires CAP_NET_BIND_SERVICE)

[Service]
User=tranzfer
Group=tranzfer
WorkingDirectory=/opt/tranzfer-mft
ExecStart=/usr/bin/java \
    -Xms256m -Xmx512m \
    -XX:+UseZGC -XX:+ZGenerational \
    -XX:+EnableDynamicAgentLoading \
    --add-opens java.base/java.lang=ALL-UNNAMED \
    --add-opens java.base/java.lang.reflect=ALL-UNNAMED \
    --add-opens java.base/java.util=ALL-UNNAMED \
    --add-opens java.base/java.lang.invoke=ALL-UNNAMED \
    --add-opens java.base/java.io=ALL-UNNAMED \
    -Dserver.port=8082 \
    -Dspring.profiles.active=production \
    -Dspring.config.additional-location=/opt/tranzfer-mft/config/ \
    -jar /opt/tranzfer-mft/lib/ftp-service.jar
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tranzfer-ftp

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/data/storage /var/log/tranzfer-mft /data/partners /tmp
PrivateTmp=true
# FTP port 21 binding: use AmbientCapabilities if needed
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
